# teave-sv/teave-sv.service
[Unit]
Description=Teave-sv systemd service.

[Service]
Type=exec
ExecStart=/usr/lib/teave/teave-sv